# Task04_custom_image_with_Qt6/meta-rpi-3b-plus/recipes-connectivity/wifi-config/files/wifi-setup.service
[Unit]
Description=WiFi Setup for wlan0
After=sys-subsystem-net-devices-wlan0.device
Wants=sys-subsystem-net-devices-wlan0.device
Before=systemd-networkd.service

[Service]
Type=simple
ExecStart=/usr/sbin/wpa_supplicant -i wlan0 -c /etc/wpa_supplicant/wpa_supplicant-wlan0.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target